# deploy/legalmind-backup.service
[Unit]
Description=LegalMind daily backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/LegalMind/deploy/backup.sh
User=root
Group=root
