# deploy/legalmind-backup.timer
[Unit]
Description=Run LegalMind backup daily

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
